# ansible/roles/camphoric/templates/camphoric.socket.j2
[Unit]
Description=camphoric socket

[Socket]
ListenStream=/run/camphoric.sock

[Install]
WantedBy=camphoric.target
