# ansible/roles/camphoric/templates/react_dev.socket.j2
[Unit]
Description=react_dev socket

[Socket]
ListenStream=/run/react_dev.sock

[Install]
WantedBy=react_dev.target
